# dockerfile/daemon/etc/s6-overlay/s6-rc.d/init-mcsm-daemon-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Auto-detect host Docker GID from mounted socket if not explicitly provided
if [ -z "${DOCKER_GID:-}" ] && [ -S /var/run/docker.sock ]; then
    DOCKER_GID=$(stat -c '%g' /var/run/docker.sock 2>/dev/null || echo "")
fi

DOCKER_GID=${DOCKER_GID:-999}

# create docker group with specified GID and add abc user to it
if getent group "$DOCKER_GID" >/dev/null 2>&1; then
    EXISTING_GROUP=$(getent group "$DOCKER_GID" | cut -d: -f1)
    echo "GID $DOCKER_GID already in use by group: $EXISTING_GROUP, adding abc to it"
    addgroup abc "$EXISTING_GROUP"
elif getent group docker >/dev/null 2>&1; then
    echo "docker group exists, adding abc to it"
    addgroup abc docker
else
    echo "Creating docker group with GID=$DOCKER_GID"
    addgroup -g "$DOCKER_GID" docker
    addgroup abc docker
fi

# permissions
lsiown -R abc:abc \
    /opt/mcsmanager/daemon

# install
cd /opt/mcsmanager/daemon && s6-setuidgid abc npm install --production
