# dockerfile/web/etc/s6-overlay/s6-rc.d/init-mcsm-web-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# permissions
lsiown -R abc:abc \
    /opt/mcsmanager/web

# install
cd /opt/mcsmanager/web && s6-setuidgid abc npm install --production
